# devel/etc/catkin/profile.d/25.turtlebot.sh
/home/vinh/catkin2_ws/devel/.private/turtlebot_bringup/etc/catkin/profile.d/25.turtlebot.sh